# dot_config/systemd/user/swayidle.service
[Unit]
Description=Swayidle Daemon
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=swayidle -w \
  timeout 300 'swaylock -f' \
  timeout 600 'swaymsg "output * power off"' \
  timeout 900 'systemctl suspend' \
  resume 'swaymsg "output * power on"' \
  before-sleep 'swaylock -f'
Restart=always

[Install]
WantedBy=graphical-session.target

