# dot_config/systemd/user/waybar.service.d/override.conf
[Unit]
ConditionEnvironment=XDG_SESSION_DESKTOP=sway
